# chrony/rootfs/etc/s6-overlay/s6-rc.d/init-chrony/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: chrony
# Configures chrony
# ==============================================================================
readonly CHRONY_CONF='/etc/chrony/chrony.conf'
declare ntp_pool
declare pool_maxsources
declare -a serverlist

# Check if configuration is valid.
if bashio::config.is_empty 'ntp_server' \
    && bashio::config.is_empty 'ntp_pool';
then
    bashio::log.fatal
    bashio::log.fatal 'Configuration of this add-on is incomplete.'
    bashio::log.fatal
    bashio::log.fatal 'Both "ntp_pool" and "ntp_server" configurations'
    bashio::log.fatal 'are empty in the add-on configuration.'
    bashio::log.fatal
    bashio::exit.nok
fi

# Write the configuration file.

ntp_pool=$(bashio::config "ntp_pool")
if [[ -n "${ntp_pool}" ]]
then
    pool_maxsources=$(bashio::config "pool_maxsources")
    bashio::log.debug "Adding pool ${ntp_pool}"
    echo "pool ${ntp_pool} iburst maxsources ${pool_maxsources}" >> ${CHRONY_CONF}
    serverlist+=("${ntp_pool}")
fi

for server in $(bashio::config "ntp_server"); do
    bashio::log.debug "Adding server ${server}"
    echo "server ${server} iburst" >> ${CHRONY_CONF}
    serverlist+=("${server}")
done

echo "initstepslew 10 ${serverlist[*]}" >> ${CHRONY_CONF}
