# chrony/rootfs/etc/services.d/chronyd/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: chrony
# Take down the S6 supervision tree when chrony fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "chrony crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "chrony stopped, restarting..."
