# chrony/rootfs/etc/services.d/chronyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: chrony
# Runs chrony daemon
# ==============================================================================
declare -a options

options+=("-d")
if bashio::config.false 'set_system_clock'
then
    options+=("-x")
fi

bashio::log.info "Starting chronyd..."

exec chronyd "${options[@]}"
