# chrony/rootfs/etc/services.d/chronyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Home Assistant Add-on: chrony
# Runs chrony daemon
# ==============================================================================
declare -a options

bashio::log.info "Starting chronyd..."

options+=(-d)
if bashio::config.false 'set_system_clock'; then
    options+=(-x)
fi

exec chronyd "${options[@]}"
